# deploy/systemd/kumomta.service.d/10-iris-diag-log.conf
# kumomta diagnostic log → file (systemd drop-in).
#
# KumoMTA writes its diagnostic/error output (ERROR/WARN/INFO/DEBUG —
# SMTP rejections, log-hook failures, policy reload notices) to stderr,
# which lands in journald. This drop-in additionally tees it to a file so
# operators have a durable, greppable log independent of journald rotation.
#
# Verbosity is controlled separately, in init.lua, via
# kumo.set_diagnostic_log_filter (rendered by iris from the
# IRIS_KUMO_DIAG_LOG_FILTER setting). This file only sets the destination.
#
# Install:
#   sudo install -Dm644 10-iris-diag-log.conf \
#     /etc/systemd/system/kumomta.service.d/10-iris-diag-log.conf
#   sudo systemctl daemon-reload && sudo systemctl restart kumomta
[Service]
StandardOutput=append:/var/log/kumomta/diagnostic.log
StandardError=append:/var/log/kumomta/diagnostic.log
