# deploy/systemd/iris-agent.service
# systemd unit for iris-agent — the KumoMTA cluster node agent.
#
# Runs on EACH remote KumoMTA node (not the central iris host). It is the same
# `iris` binary invoked with the `agent` subcommand: an mTLS server that receives
# config bundles from the iris control plane, writes kumod's generated policy,
# reloads/restarts kumod, reports health, and reverse-proxies kumod's localhost
# admin API. See docs/kumomta-cluster-architecture.md.
#
# Runs as the SAME user kumod runs as (iris in the typical `kumod --user iris`
# setup), so the policy files it writes are owned by that user and kumod can
# read them directly — no config_group/setgid needed. Because a non-root user
# cannot ask systemd to restart a service, restarts (needed only for init-block
# changes: listeners/spool/log hook) go through sudo; install the companion
# iris-agent.sudoers to allow exactly those two systemctl calls.
#
# Prerequisites on this node:
#   - the `iris` binary installed at /usr/local/bin/iris
#   - the policy dir writable by iris:  chown -R iris:iris /opt/kumomta/etc/policy
#   - /etc/iris/iris.yaml with:
#       kumomta:   stub: false, config_path, base_url,
#                  restart_command: "sudo systemctl restart kumomta"
#                  (reloads use kumod's 10s file-watcher — leave reload_* empty)
#       agent:     listen, ca_cert, cert, key  (the certs from `iris cluster enroll`)
#   - install the sudoers rule:  cp iris-agent.sudoers /etc/sudoers.d/iris-agent
#     (then: chmod 440 /etc/sudoers.d/iris-agent && visudo -c)
#   - the node registered + enrolled in iris (Cluster page -> Enroll)
#
# Install:
#   cp iris-agent.service /etc/systemd/system/iris-agent.service
#   cp iris-agent.sudoers /etc/sudoers.d/iris-agent && chmod 440 /etc/sudoers.d/iris-agent
#   systemctl daemon-reload
#   systemctl enable --now iris-agent
#   journalctl -u iris-agent -f

[Unit]
Description=iris-agent — KumoMTA cluster node agent (mTLS control plane)
Documentation=https://github.com/menta2k/iris
After=network-online.target kumomta.service
Wants=network-online.target

[Service]
Type=simple
# Run as the user kumod runs as, so policy files are readable by kumod without
# config_group/setgid. Change this if kumod runs as a different user (e.g.
# `kumod`), and update the sudoers file's user to match.
User=iris
Group=iris

# The agent needs no DB/Redis; it reads kumomta:/agent: from the yaml.
# `-` = tolerate a missing file.
EnvironmentFile=-/etc/iris/iris.env

# NOTE: -config MUST come BEFORE the `agent` subcommand (Go flag parsing stops
# at the first positional argument).
ExecStart=/usr/local/bin/iris -config /etc/iris/iris.yaml agent

Restart=on-failure
RestartSec=5s

# Hardening — the agent only needs to read its config/certs, write the kumod
# policy dir, and reach the network (its mTLS listener + localhost kumod).
# NoNewPrivileges must stay false: the agent shells out to
# `sudo systemctl restart kumomta`, and sudo needs to raise privileges.
NoNewPrivileges=false
PrivateTmp=true
ProtectHome=true

[Install]
WantedBy=multi-user.target
