# deploy/systemd/iris.service
# iris — KumoMTA operator console; systemd unit for a host-native install.
#
# Runs the iris binary alongside a host-installed kumomta on the
# same machine, with PostgreSQL/TimescaleDB and Redis reachable on
# loopback. Configuration is split:
#
#   /etc/iris/iris.env   secrets + env knobs (root-only)
#   /etc/iris/configs/            kratos layered YAML (data, server, auth)
#   /var/lib/iris/                (reserved — currently unused)
#   /opt/kumomta/etc/policy/init.lua    written by admin-service, read by kumomta
#   /opt/kumomta/etc/dkim/              DKIM PEM keys
#
# Install:
#   1. sudo cp deploy/systemd/iris.service /etc/systemd/system/
#   2. sudo cp deploy/systemd/iris.env.example /etc/iris/iris.env
#      and edit secrets
#   3. sudo cp -r backend/app/admin/service/configs /etc/iris/configs
#   4. install the binary: sudo install -m 0755 ./iris /usr/local/bin/
#   5. sudo systemctl daemon-reload && sudo systemctl enable --now iris
#
# Logs:
#   sudo journalctl -u iris -f

[Unit]
Description=iris — KumoMTA operator console (UI + admin API)
Documentation=https://github.com/menta2k/iris
# Wait for the data services AND for kumomta to be up. kumomta is on the
# same host so its unit name is conventional; adjust if your packaging
# names it differently.
After=network-online.target postgresql.service redis.service kumomta.service
Wants=network-online.target

[Service]
Type=simple
User=iris
Group=kumomta
# kumomta group membership: admin-service WRITES into the policy and
# DKIM directories that kumomta reads. Putting the binary in `kumomta`
# group makes the 0664/0640 permissions work without a world-readable
# DKIM key directory.

EnvironmentFile=/etc/iris/iris.env

# --conf points at the layered Kratos configs. Override per-knob via the
# environment file rather than editing this unit.
ExecStart=/usr/local/bin/iris --conf /etc/iris/configs

# Restart policy: admin-service is stateless (state is in PG+Redis),
# crashes are safe to restart aggressively. Backoff prevents tight loops
# on a config error.
Restart=on-failure
RestartSec=5s

# Hardening. The binary needs to read configs and write
# /opt/kumomta/etc/policy/ — both whitelisted via ReadWritePaths.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
PrivateDevices=true
PrivateTmp=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallFilter=@system-service

# Writable surfaces. /opt/kumomta/etc/{policy,dkim} are the bind-mount
# equivalents from docker-compose: admin-service generates init.lua and
# emits/imports DKIM PEMs there; kumomta reads them.
ReadWritePaths=/opt/kumomta/etc/policy /opt/kumomta/etc/dkim /var/lib/iris

# Network: loopback-only by default would require socket-level binding,
# which kratos doesn't expose. We let the binary bind whatever
# server.yaml says (typically 127.0.0.1:8000) and trust the systemd
# RestrictAddressFamilies for limiting to AF_INET.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Resource limits. NOFILE 65k matches kumomta's own; admin-service holds
# Redis pool conns, DB pool conns, and HTTP keep-alives, but nothing
# fundamentally fan-out heavy.
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
