# deploy/nfpm/iris.service
# systemd unit for the Iris admin backend.
#
# Installed to /usr/lib/systemd/system/iris.service by the deb/rpm package.
# Enable + start with:
#   systemctl enable --now iris
#
# Secrets and overrides live in /etc/iris/iris.env (created from
# iris.env.example on first install). The non-secret config is
# /etc/iris/iris.yaml.

[Unit]
Description=Iris — KumoMTA operator console
Documentation=https://github.com/menta2k/iris
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=iris
Group=iris

# IRIS_* variables here override values in iris.yaml (see conf.applyEnv):
# IRIS_DATABASE_DSN, IRIS_REDIS_ADDR, IRIS_REDIS_PASSWORD,
# IRIS_SESSION_SECRET, IRIS_HTTP_ADDR, IRIS_GRPC_ADDR, IRIS_LOG_LEVEL, etc.
EnvironmentFile=/etc/iris/iris.env

ExecStart=/usr/local/bin/iris -config /etc/iris/iris.yaml

Restart=on-failure
RestartSec=5s

# Hardening — the service only needs to read its config and write its own
# state directory.
StateDirectory=iris
WorkingDirectory=/var/lib/iris
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/iris

# Iris talks to PostgreSQL and Redis over TCP, and Kratos enumerates network
# interfaces (via a netlink socket) to pick its advertised endpoint. Some
# hardened hosts set a restrictive DefaultRestrictAddressFamilies in
# system.conf, which would make socket(AF_INET/AF_NETLINK, ...) fail with
# EAFNOSUPPORT ("address family not supported by protocol") for any unit that
# doesn't opt in. Declare the families we need so the service works regardless
# of the host default:
#   AF_UNIX    — local sockets (journal, etc.)
#   AF_INET(6) — TCP to PostgreSQL / Redis and the HTTP/gRPC listeners
#   AF_NETLINK — interface enumeration during server startup
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

[Install]
WantedBy=multi-user.target
