# deploy/nfpm/iris.service
# systemd unit for the Iris admin backend.
#
# Installed to /usr/lib/systemd/system/iris.service by the deb/rpm package.
# Enable + start with:
#   systemctl enable --now iris
#
# Secrets and overrides live in /etc/iris/iris.env (created from
# iris.env.example on first install). The non-secret config is
# /etc/iris/iris.yaml.

[Unit]
Description=Iris — KumoMTA operator console
Documentation=https://github.com/menta2k/iris
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=iris
Group=iris

# IRIS_* variables here override values in iris.yaml (see conf.applyEnv):
# IRIS_DATABASE_DSN, IRIS_REDIS_ADDR, IRIS_REDIS_PASSWORD,
# IRIS_SESSION_SECRET, IRIS_HTTP_ADDR, IRIS_GRPC_ADDR, IRIS_LOG_LEVEL, etc.
EnvironmentFile=/etc/iris/iris.env

ExecStart=/usr/local/bin/iris -config /etc/iris/iris.yaml

Restart=on-failure
RestartSec=5s

# Allow binding privileged ports (< 1024) as the non-root iris user, so an
# admin or injection listener can be configured on e.g. :443 or :80 without
# running as root. Applied by systemd at exec. Remove both lines if you always
# front Iris with a reverse proxy and never bind low ports directly.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Hardening — the service only needs to read its config and write its own
# state directory.
StateDirectory=iris
WorkingDirectory=/var/lib/iris
# NoNewPrivileges must stay false on a co-located host (iris + kumod on the same
# machine): applying an init-block config change restarts kumod via
# kumomta.restart_command, which is "sudo systemctl restart kumomta" (iris runs
# as the non-root iris user). NoNewPrivileges=true would make sudo fail with
# "the 'no new privileges' flag is set, which prevents sudo from running as
# root". The companion iris-agent.sudoers scopes that sudo to exactly the
# restart/reload of kumomta. If this host is control-plane ONLY (every kumod is
# remote, reached via an iris-agent) you can tighten this back to true with a
# drop-in: /etc/systemd/system/iris.service.d/*.conf → [Service] NoNewPrivileges=true
NoNewPrivileges=false
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/iris

# Iris talks to PostgreSQL and Redis over TCP, and Kratos enumerates network
# interfaces (via a netlink socket) to pick its advertised endpoint. Some
# hardened hosts set a restrictive DefaultRestrictAddressFamilies in
# system.conf, which would make socket(AF_INET/AF_NETLINK, ...) fail with
# EAFNOSUPPORT ("address family not supported by protocol") for any unit that
# doesn't opt in. Declare the families we need so the service works regardless
# of the host default:
#   AF_UNIX    — local sockets (journal, etc.)
#   AF_INET(6) — TCP to PostgreSQL / Redis and the HTTP/gRPC listeners
#   AF_NETLINK — interface enumeration during server startup
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

[Install]
WantedBy=multi-user.target
